# app/supervisord.d/httpd.conf
[program:httpd]
command=/usr/sbin/apachectl -DFOREGROUND
autostart=true
autorestart=true
priority=5
redirect_stderr=true
stdout_logfile=/var/log/httpd.log
