# app/supervisord.d/monitord.conf
[program:monitord]
command=/app/bin/start-monitord
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/monitord.log
