# app/supervisord.d/netbootmond.conf
[program:netbootmond]
command=/app/bin/start-netbootmond
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/netbootmond.log
