# app/supervisord.d/php-fpm.conf
[program:php-fpm]
command=/usr/sbin/php-fpm -F
autostart=true
autorestart=true
priority=5
redirect_stderr=true
stdout_logfile=/var/log/php-fpm.log
