# app/supervisord.d/rsyncd.conf
[program:rsyncd]
command=rsync --daemon --no-detach --config /etc/systemimager/rsyncd.conf
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/rsyncd.log
