# app/supervisord.d/tftpd.conf
[program:tftpd]
command=/usr/sbin/in.tftpd --foreground -vvv -s /var/lib/tftpboot
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tftp.log
